# com.user.codereview.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ilya.codereview</string>

    <key>ProgramArguments</key>
    <array>
        <string>__VENV__/bin/python</string>
        <string>__AGENT_DIR__/agent.py</string>
        <string>--repos-root</string>
        <string>__REPOS_ROOT__</string>
    </array>

    <!-- Run once daily Mon–Fri at 23:00 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>23</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>23</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>23</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>23</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>23</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/review-agent.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/review-agent.log</string>

    <!-- Restart on crash, but not immediately -->
    <key>ThrottleInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
